# supervisord.conf
[supervisord]
nodaemon=true

[program:shellinabox]
command=/usr/bin/shellinaboxd --disable-ssl
numprocs=1
autostart=true
autorestart=true
