# supervisord.conf
[supervisord]
nodaemon=true

[program:shellinabox]
command=/usr/bin/shellinaboxd --disable-ssl --css=/etc/shellinabox/shellinabox.css
numprocs=1
autostart=true
autorestart=true
